# skeybindd.service
[Unit]
Description=A TTY Keybind Daemon

[Service]
ExecStart=/usr/bin/skeybindd

[Install]
WantedBy=default.target