# skeybindd.service
[Unit]
Description=TTY Keybind Daemon
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/skeybindd

[Install]
WantedBy=default.target
